# grid/ec2/services/manifest/aura-create-nondataservices.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type='manifest' name='aura:import-data'>

<service name='aura/create-nondataservices' type='service' version='1'>
	<create_default_instance enabled='true' />
	<single_instance />
        <dependency name='multi-user' grouping='require_all' restart_on='none' type='service'>
    	    <service_fmri value='svc:/milestone/multi-user-server:default' />
	</dependency>
	<exec_method type='method' name='start' 
		     exec='/opt/aura/services/method/createServices.sh no' timeout_seconds='320' />
	<exec_method type='method' name='stop' exec=':true' timeout_seconds='0' />	
	<exec_method type='method' name='refresh' exec=':true' timeout_seconds='0' />	
	<exec_method type='method' name='restart' exec=':true' timeout_seconds='0' />	
	<property_group name='startd' type='framework'>
		<propval name='duration' type='astring' value='transient' />
	</property_group>
	<stability value='Unstable' />
	<template>
		<common_name>
			<loctext xml:lang='C'>Create the services to run Aura components that don't depend on having a datapool mounted.</loctext>
		</common_name>
	</template>
</service>
</service_bundle>
